# 01-xinput-disable.rules
RUN+="/bin/sh -c '/usr/bin/xinput list --name-only | grep ^xwayland-pointer-gestures | /usr/bin/xargs -n1 /usr/bin/xinput disable'"
